# install/linux/avadhi@.service
[Unit]
Description=Avadhi Time Collector (%i)
Documentation=https://github.com/urwithajit9/avadhi-collector
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot

# Dedicated service user
User=avadhi

# Ensure relative paths resolve correctly
WorkingDirectory=/opt/avadhi-collector

# Main execution
ExecStart=/opt/avadhi-collector/avadhi-collector run

# Hard preflight checks
ExecStartPre=/usr/bin/test -x /opt/avadhi-collector/avadhi-collector
ExecStartPre=/usr/bin/test -f /opt/avadhi-collector/Config.toml
ExecStartPre=/usr/bin/test -f /opt/avadhi-collector/AvadhiConfig.toml

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
# IMPORTANT:
# This service MUST NOT be enabled directly.
# It is triggered only by avadhi.timer.
WantedBy=timers.target
