# install/linux/avadhi.service
[Unit]
Description=Avadhi Collector Dispatcher

[Service]
Type=oneshot
ExecStart=/bin/systemctl start avadhi@%u.service
