# install/linux/avadhi.service
[Unit]
Description=Avadhi Time Collector Service
Documentation=https://github.com/urwithajit9/avadhi-collector
After=network.target

[Service]

The service runs as the user specified in the instance name (e.g., avadhi@ajit.service runs as user 'ajit').

User=%i
Group=%i

=== CRITICAL FIX: Working Directory (Resolves Config.toml not found in service) ===

Set the current working directory to the installation path, ensuring the binary can find

Config.toml and AvadhiConfig.toml.

WorkingDirectory=/opt/avadhi-collector

Execute the binary from the fixed installation path

ExecStart=/opt/avadhi-collector/avadhi-collector

Restart the service on failure for resilience

Restart=always
RestartSec=5s

Standard output and error logging

StandardOutput=journal
StandardError=journal

[Install]

This enables the service to be started when the system enters the multi-user state

WantedBy=multi-user.target